# server-audit/backup/systemd/backup-home.timer
[Unit]
Description=Daily Incremental Home Backup to MEGA (rclone)

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=15m

[Install]
WantedBy=timers.target
